# spot-hinta-mcp.service
[Unit]
Description=spot-hinta MCP server (Finnish electricity spot prices)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=username
WorkingDirectory=/home/username/gitrepos/spot-mcp-server
ExecStart=/home/username/gitrepos/spot-mcp-server/venv/bin/python /home/username/gitrepos/spot-mcp-server/spot_hinta_mcp.py
Environment=SPOT_HINTA_PORT=8767
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
