# spot-hinta-mcp.service
[Unit]
Description=spot-hinta MCP server (Finnish electricity spot prices)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jsavola
WorkingDirectory=/home/jsavola/gitrepos/spot-mcp-server
ExecStart=/usr/bin/python3 /home/jsavola/gitrepos/spot-mcp-server/spot_hinta_mcp.py
Environment=SPOT_HINTA_PORT=8767
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
